# web/heap-dump/heap-dump/rootfs/etc/supervisor/conf.d/secretinit.conf
[program:secretinit]
command=/secretinit.sh
user=root
startretries=0
startsecs=0
autorestart=false
priority=3